# usbip/server/99-usbip-autobind.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/local/bin/usbip-autobind.sh"
