# audio-streamer/WWH23-feed.service
[Unit]
Description=WWH23 Weather Radio Stream
After=network-online.target icecast2.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/WWH23-feed.sh
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
